# software/winglet-os/winglet/overlay/etc/profile.d/qt_env.sh
export QT_FONT_DPI=84
export QT_QPA_FB_DRM=1
export QT_QPA_PLATFORM=linuxfb:nographicsmodeswitch
export QT_QPA_GENERIC_PLUGINS=evdevmouse:/dev/inputs-by-name/rotary_encoder,evdevkeyboard:/dev/inputs-by-name/gpio_buttons,evdevkeyboard:/dev/inputs-by-name/axp2585-pek
