# software/winglet-os/winglet/overlay/etc/profile.d/shell_prompt.sh
PS1='\[\e[90m\][\u@\h]\[\e[0m\] \[\e[36m\]\W\[\e[0m\] \$ '
