# software/winglet-os/winglet/overlay/etc/network/if-up.d/udhcpd-usb0
#!/bin/sh

if ! [ "$IFACE" = "usb0" ]; then
    # Not our interface, quit
    exit 0
fi

udhcpd /etc/udhcpd.conf
exit $?
